# packer/files/systemd/openclaw.service
[Unit]
Description=OpenClaw Enterprise AI Agent
Documentation=https://github.com/openclaw/openclaw
After=network.target docker.service litellm.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/openclaw/openclaw
EnvironmentFile=-/opt/openclaw/openclaw.env
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=30

ExecStartPre=/bin/bash -c 'echo "Starting OpenClaw..."'

MemoryHigh=80%
MemoryMax=4G
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
